# deploy/serabutd.service
[Unit]
Description=Serabut PXE Boot Manager
Documentation=https://github.com/twdamhore/serabut
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/serabutd
Restart=on-failure
RestartSec=5

# Environment
Environment=RUST_LOG=serabutd=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Allow read access to data directory
ReadWritePaths=/var/lib/serabutd

# Allow binding to port 8080 (no special privileges needed for ports > 1024)
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
