# deploy/serabutd.service
[Unit]
Description=Serabut PXE Boot Server
Documentation=https://github.com/twdamhore/serabut
After=network.target

[Service]
Type=simple
User=serabut
Group=serabut
ExecStart=/usr/local/bin/serabutd
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/serabutd
PrivateTmp=yes
LogsDirectory=serabut

[Install]
WantedBy=multi-user.target
